# hello.service
[Unit]
Description=Hello World Service
After=basic.target

[Service]
User=root
Type=simple
ExecStart=/usr/bin/hello
ExecStopPost=/usr/bin/kill $MAINPID
RemainAfterExit=no
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
